# debian/postinst
#!/bin/sh
# postinst script for dbus-mqtt-bridge
# SPDX-License-Identifier: GPL-3.0-or-later
# Copyright (C) 2026 Ed Lee

set -e

# Summary of what this script does:
# 1. Creates system user and group
# 2. Creates necessary directories with correct permissions
# 3. Installs example config if none exists
# 4. Sets up D-Bus policy template
# 5. Displays configuration instructions

case "$1" in
    configure)
        # Create system user and group if they don't exist
        if ! getent group dbus-mqtt-bridge >/dev/null 2>&1; then
            addgroup --system dbus-mqtt-bridge
        fi
        
        if ! getent passwd dbus-mqtt-bridge >/dev/null 2>&1; then
            adduser --system \
                    --ingroup dbus-mqtt-bridge \
                    --home /nonexistent \
                    --no-create-home \
                    --gecos "D-Bus MQTT Bridge Service" \
                    --shell /usr/sbin/nologin \
                    dbus-mqtt-bridge
        fi
        
        # Create config directory
        if [ ! -d /etc/dbus-mqtt-bridge ]; then
            mkdir -p /etc/dbus-mqtt-bridge
        fi
        
        # Create log directory
        if [ ! -d /var/log/dbus-mqtt-bridge ]; then
            mkdir -p /var/log/dbus-mqtt-bridge
            chown dbus-mqtt-bridge:dbus-mqtt-bridge /var/log/dbus-mqtt-bridge
            chmod 0750 /var/log/dbus-mqtt-bridge
        fi
        
        # Install example config if no config exists
        if [ ! -f /etc/dbus-mqtt-bridge/config.yaml ]; then
            if [ -f /usr/share/doc/dbus-mqtt-bridge/examples/config.yaml ]; then
                # Handle compressed example file
                if [ -f /usr/share/doc/dbus-mqtt-bridge/examples/config.yaml.gz ]; then
                    zcat /usr/share/doc/dbus-mqtt-bridge/examples/config.yaml.gz > /etc/dbus-mqtt-bridge/config.yaml
                else
                    cp /usr/share/doc/dbus-mqtt-bridge/examples/config.yaml /etc/dbus-mqtt-bridge/config.yaml
                fi
                
                # Set ownership and permissions
                chown root:dbus-mqtt-bridge /etc/dbus-mqtt-bridge/config.yaml
                chmod 0640 /etc/dbus-mqtt-bridge/config.yaml
                
                echo "Installed example configuration to /etc/dbus-mqtt-bridge/config.yaml"
            fi
        else
            # Ensure correct permissions on existing config
            chown root:dbus-mqtt-bridge /etc/dbus-mqtt-bridge/config.yaml 2>/dev/null || true
            chmod 0640 /etc/dbus-mqtt-bridge/config.yaml 2>/dev/null || true
        fi
        
        # Set directory permissions
        chown root:dbus-mqtt-bridge /etc/dbus-mqtt-bridge
        chmod 0750 /etc/dbus-mqtt-bridge
        
        # Install D-Bus policy template if none exists
        # Note: The actual policy must be customized by the user
        if [ ! -f /etc/dbus-1/system.d/dbus-mqtt-bridge.conf ]; then
            if [ -f /usr/share/doc/dbus-mqtt-bridge/examples/dbus-policy.conf ]; then
                if [ -f /usr/share/doc/dbus-mqtt-bridge/examples/dbus-policy.conf.gz ]; then
                    zcat /usr/share/doc/dbus-mqtt-bridge/examples/dbus-policy.conf.gz > /etc/dbus-1/system.d/dbus-mqtt-bridge.conf
                else
                    cp /usr/share/doc/dbus-mqtt-bridge/examples/dbus-policy.conf /etc/dbus-1/system.d/dbus-mqtt-bridge.conf
                fi
                chmod 0644 /etc/dbus-1/system.d/dbus-mqtt-bridge.conf
                
                echo "Installed D-Bus policy template to /etc/dbus-1/system.d/dbus-mqtt-bridge.conf"
            fi
        fi
        
        # Reload D-Bus configuration
        if [ -x /usr/bin/dbus-send ]; then
            dbus-send --system --type=method_call --dest=org.freedesktop.DBus \
                     / org.freedesktop.DBus.ReloadConfig >/dev/null 2>&1 || true
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# Display configuration instructions (only on fresh install)
if [ "$1" = "configure" ] && [ -z "$2" ]; then
    cat <<EOF

================================================================================
dbus-mqtt-bridge has been installed but NOT enabled.

IMPORTANT CONFIGURATION STEPS:

1. Edit the configuration file:
   sudo editor /etc/dbus-mqtt-bridge/config.yaml

2. Customize the D-Bus policy for your specific mappings:
   sudo editor /etc/dbus-1/system.d/dbus-mqtt-bridge.conf
   
   The installed template includes examples. You MUST modify it to match
   the D-Bus services, interfaces, and methods in your config.yaml.

3. Validate your configuration:
   dbus-mqtt-bridge --help    # See command-line options
   
4. After configuration, enable and start the service:
   sudo systemctl enable --now dbus-mqtt-bridge

SECURITY NOTES:
- The service runs as user 'dbus-mqtt-bridge' (not root)
- D-Bus access is restricted by /etc/dbus-1/system.d/dbus-mqtt-bridge.conf
- Config file is only readable by root and dbus-mqtt-bridge group

For detailed information:
   man dbus-mqtt-bridge

================================================================================
EOF
fi

exit 0